# frontend/root/etc/services.d/frontend/run
#!/command/with-contenv bash

set 

exec s6-setuidgid abc \
  caddy /default/Caddyfile
